# provision/modules/sentry/templates/sentry.service.erb
#
# This file is generated by Puppet, don't touch please
# sentry::install
#
# module: sentry
# source: sentry-celery.service.erb
#


[Unit]
Description=Celery background worker
After=network.target multi-user.target

[Service]
User=<%= scope['sentry::params::user'] %>
Group=<%= scope['sentry::params::group'] %>
Environment=SENTRY_CONF=/etc/sentry/
WorkingDirectory=/opt/sentry

ExecStart=/usr/bin/sentry start
ExecReload=/bin/kill -HUP $MAINPID

