# provision/modules/sentry/templates/sentry-celery-beat.service.erb
#
# This file is generated by Puppet, don't touch please
#
# module: sentry
# source: sentry-celery-beat.service.erb
#

[Unit]
Description=Celery beat worker
After=network.target multi-user.target

[Service]
User=<%= scope['sentry::user'] %>
Group=<%= scope['sentry::group'] %>
Environment=SENTRY_CONF=/etc/sentry/
WorkingDirectory=/opt/sentry

ExecStart=/usr/bin/sentry celery beat
ExecReload=/bin/kill -HUP $MAINPID

